# automate_infra/report_monitor.service
[Unit]
Description=Monitor CAPE reports for categorization
After=network.target

[Service]
ExecStart=/usr/bin/python3 /home/cape/scripts/categorize_reports.py
Restart=always
User=cape

[Install]
WantedBy=multi-user.target
